# deploy/services/gddo-server/run
#!/bin/sh
exec /go/bin/gddo-server -robot=999999999 -trust_proxy_headers -sourcegraph_url=https://sourcegraph.com 2>&1
