# deploy/services/gddo-server/run
#!/bin/sh
exec /go/bin/gddo-server -trust_proxy_headers -sourcegraph_url=https://staging2.sourcegraph.com 2>&1
